# automation/claude-research-bot.service
[Unit]
Description=claude-research-bot ARO — one autonomous research cycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tar/claude-research-bot
ExecStart=/home/tar/claude-research-bot/scripts/supervisor.sh
# claude CLI lives in ~/.local/bin
Environment=HOME=/home/tar
Environment=PATH=/home/tar/.local/bin:/usr/local/bin:/usr/bin:/bin
# A cycle should not run forever; supervisor also enforces its own timeout.
TimeoutStartSec=3600
Nice=10
# Keep full-auto runs scoped to the research bot and local Claude state.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/tar/claude-research-bot /home/tar/.claude /home/tar/.config /home/tar/.local /tmp
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=default.target
